# wordcloud.timer
[Unit]
Description=Post wordcloud every hour.

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
